# overlay/common/etc/init.d/hcos-agent
#!/sbin/openrc-run

name="Home Cloud OS agent"
command="/usr/sbin/hcos"
command_args="agent"
